# ops/systemd/rbotzilla-backend.service
[Unit]
Description=RBOTzilla Backend Service
After=network.target

[Service]
Type=simple
User=%u
WorkingDirectory=/home/%u/RICK_PHOENIX
Environment=PYTHONUNBUFFERED=1
Environment=RICK_ENV=practice
Environment=ALLOW_PRACTICE_ORDERS=1
Environment=CONFIRM_PRACTICE_ORDER=1
Environment=PRACTICE_PIN=841921
ExecStart=/home/%u/.venv/bin/python /home/%u/RICK_PHOENIX/backend.py
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
