# ROLLBACK_SNAPSHOTS/live_backup_1764825818/ops/rick_phoenix.service
[Unit]
Description=RICK Phoenix System - Rick Agent Launcher
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/app
ExecStart=/usr/bin/bash /app/start_with_integrity.sh
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
Environment=ENV=production

[Install]
WantedBy=multi-user.target
